# alicexrdservices.tmp
#!/bin/bash
#       /etc/rc.d/init.d/xrd
# Starts xrootd
#
# chkconfig: 345 95 5
# description: Starts xrootd
# processname: xrootd

# Source function library.

if test ! -x /usr/bin/xrootd ; then
  echo Error: /usr/bin/xrootd does not exist
  exit 0
fi


if test ! -x XRDSHDIR/xrd.sh ; then
  echo Error: XRDSHDIR/xrd.sh does not exist
  exit 0
fi

start() {
        echo -n "Starting xrootd services: "
        XRDSHDIR/xrd.sh -c
        return $?
}

stop() {
        echo -n "Stopping xrootd services: "
        XRDSHDIR/xrd.sh -k
        RETVAL=$?
        [ $RETVAL -eq 0 ] && rm -f /var/lock/subsys/xrootd
        echo
        return $RETVAL
}


restart() {
        stop
        start
}

reload() {
        restart
}


case "$1" in
start)
        start
		RETVAL=$?
        ;;
stop)
        stop
		RETVAL=$?
        ;;
reload|restart)
        restart
		RETVAL=$?
        ;;
condrestart)
        if [ -f /var/lock/subsys/xrootd ]; then
            restart
			RETVAL=$?
		else
			RETVAL=0
        fi
        ;;
status)
		XRDSHDIR/xrd.sh -s
		RETVAL=$?
        ;;
*)
		echo "Usage: $0 {start|stop|restart|condrestart|status}"
        exit 1
esac

exit $RETVAL
